# etc/supervisor/services/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -c /opt/mailman3/etc/nginx/nginx.conf
process_name=nginx
priority=2
redirect_stderr=true    
stdout_logfile=/opt/mailman3/var/logs/nginx-stdout.log
