# etc/supervisor/services/webui.ini
[program:webui]
command=/opt/mailman3/conda/envs/mailman3_ext/bin/waitress-serve --listen=localhost:8787 webui.wsgi:application
process_name=webui
priority=1
directory=/opt/mailman3/webui
redirect_stderr=true          
stdout_logfile=/opt/mailman3/var/logs/webui-stdout.log
